# profile.d/fzf.sh
eval "$(fzf --bash)"
alias f="fzf --multi --style=full --preview='fzf-preview.sh {}'"
# alias f="fzf --tmux --preview='bat --style=numbers --color=always {}'"
