# profile.d/git.sh
alias gc="git commit"
alias gs="git status"
alias gb="git branch"
alias gba="git branch --all"
alias gdc="git diff --color"
alias gdch="git diff --color HEAD"
alias gdcc="git diff --word-diff=color"
