# profile.d/git.sh
alias gs="git status"
alias gb="git branch"
alias gdc="git diff --color"
alias gdcc="git diff --word-diff=color"
