# profile.d/kubectl.sh
alias k="kubectl"
alias ka="kubectl apply -f"
alias kd="kubectl describe"
alias kdp="kubectl describe pod"
